# php-project/project-of-owner/queue_config/quotesStatusUpdate.ini
[program:quotesStatusUpdate]#指定当前文件名必须唯一
process_name=%(program_name)s_%(process_num)02d
command= php /alidata/www/fiberstore.com/app/Command/QuotesStatusUpdateCommand.php
autostart=true ;在supervisord启动时自动启动，默认为true
autorestart=true ;是否自动重启，可选参数为 false, unexpected, true。如果为false则不自动重启，如果为unexpected表示如果程序退出信号不在 `exitcodes` 中，则自动重启。默认为unexpected
user=www ;电脑用户名
numprocs=1 ;程序运行的副本个数，默认为1，如果值大于1，则`process_name` 必须包含 `%(process_num)s`
redirect_stderr=true ;将错误输出定向到标准输出，默认为false
stdout_logfile=/data/logs/supervisor-pro/quotesStatusUpdate.log  ;program_name.log 为应用程序log
stopwaitsecs=3600 ;发送SIGKILL信号前最大等待时间。默认为10。
startretries=5 ;尝试重启最大次数。默认为3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
